# crt-daemon.service
[Unit]
Description=Apple Studio Display CRT Daemon
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/crt-control/crt-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
