# crt-daemon.service
[Unit]
Description=Apple Studio Display CRT Power-On
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/lib/crt-control/crt-daemon
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
