# debian/postinst
#!/bin/bash
set -e

# Enable and start the daemon
systemctl daemon-reload
systemctl enable crt-daemon.service
systemctl start crt-daemon.service || true

# Install labwc keybinding (Super+M → power toggle)
LABWC_DIR="/etc/xdg/labwc"
RC_FILE="$LABWC_DIR/rc.xml"

if [ -d "$LABWC_DIR" ]; then
    if [ -f "$RC_FILE" ]; then
        # Check if keybinding already exists
        if ! grep -q "crt-toggle" "$RC_FILE" 2>/dev/null; then
            # Insert keybinding before closing </openbox_config> tag
            sed -i 's|</openbox_config>|  <keyboard>\n    <keybind key="C-A-F12">\n      <action name="Execute"><command>sudo /usr/bin/crt-toggle</command></action>\n    </keybind>\n  </keyboard>\n</openbox_config>|' "$RC_FILE"
        fi
    else
        cp /usr/share/crt-control/99-crt-control.xml "$RC_FILE"
    fi
fi

# Add sudoers rule for passwordless crt-toggle and GUI
cat > /etc/sudoers.d/crt-control << 'EOF'
ALL ALL=(root) NOPASSWD: /usr/bin/crt-toggle
ALL ALL=(root) NOPASSWD: /usr/bin/python3 /usr/lib/crt-control/pi-adc-gui.py
EOF
chmod 440 /etc/sudoers.d/crt-control

echo "pi-adc-control installed. Super+M toggles monitor power."
echo "Open Preferences → ADC CRT Control for geometry adjustments."
