# debian/prerm
#!/bin/bash
set -e

# Stop and disable daemon
systemctl stop crt-daemon.service || true
systemctl disable crt-daemon.service || true

# Remove keybinding from labwc
RC_FILE="/etc/xdg/labwc/rc.xml"
if [ -f "$RC_FILE" ]; then
    # Remove the crt-toggle keybinding block
    sed -i '/crt-toggle/d' "$RC_FILE"
    sed -i '/Super_L-m/d' "$RC_FILE"
fi

# Remove sudoers rule
rm -f /etc/sudoers.d/crt-control

systemctl daemon-reload
